# migration_sources/omniarchon/deployment/archon-alerting.service
# Archon Slack Alerting Systemd Service
#
# Installation:
# 1. Create webhook credential file:
#    echo "https://hooks.slack.com/services/YOUR/WEBHOOK/URL" | sudo tee /etc/archon/slack_webhook.txt
#    sudo chmod 600 /etc/archon/slack_webhook.txt
#    sudo chown root:root /etc/archon/slack_webhook.txt
# 2. Edit this file and update paths and user
# 3. Copy to systemd directory:
#    sudo cp deployment/archon-alerting.service /etc/systemd/system/
# 4. Reload systemd:
#    sudo systemctl daemon-reload
# 5. Enable and start:
#    sudo systemctl enable archon-alerting
#    sudo systemctl start archon-alerting
# 6. Check status:
#    sudo systemctl status archon-alerting
# 7. View logs:
#    sudo journalctl -u archon-alerting -f

[Unit]
Description=Archon Slack Alerting Service
Documentation=https://github.com/OmniNode-ai/omniarchon/blob/main/docs/SLACK_ALERTING.md
After=docker.service
Requires=docker.service
# Optional: Wait for network if Docker socket is remote
# After=network-online.target
# Wants=network-online.target

[Service]
Type=simple
# Change to your user (must have Docker access)
User=your-user
Group=docker
# Change to your repository path
# PLACEHOLDER: Update this to your actual repository location
# Options:
#   - Absolute path: /home/youruser/omniarchon
#   - Home-relative: Use %h to reference home directory (e.g., %h/omniarchon)
#   - Environment variable: Use ${MY_REPO_PATH} and set via Environment= directive
# Example below uses development machine path - MUST BE CHANGED for deployment
WorkingDirectory=/Volumes/PRO-G40/Code/omniarchon

# Load webhook URL securely using systemd credentials
# The webhook URL will be available at: $CREDENTIALS_DIRECTORY/slack_webhook
# This prevents the URL from being exposed via /proc/<pid>/environ
# REQUIRED: Create credential file at /etc/archon/slack_webhook.txt (see installation steps above)
LoadCredential=slack_webhook:/etc/archon/slack_webhook.txt

# Optional: Override default settings via environment variables
# Environment="MONITORING_CHECK_INTERVAL_SECONDS=60"
# Environment="ALERT_THRESHOLD_CPU_PERCENT_WARNING=70.0"
# Environment="ALERT_THROTTLE_MAX_ALERTS_PER_WINDOW=2"

# Note: The Python script should read the webhook URL from:
#   webhook_path = os.environ.get('CREDENTIALS_DIRECTORY', '') + '/slack_webhook'
#   with open(webhook_path, 'r') as f:
#       webhook_url = f.read().strip()
# Falls back to ALERT_NOTIFICATION_SLACK_WEBHOOK_URL environment variable for backward compatibility
# Use full path to Python (find with: which python3)
ExecStart=/usr/bin/python3 scripts/slack_alerting.py --daemon
# Restart on failure
Restart=always
RestartSec=10
# Standard output to journal
StandardOutput=journal
StandardError=journal
# Security hardening (optional)
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ReadWritePaths=/tmp  # For state file

[Install]
WantedBy=multi-user.target
